# etc/systemd/system/etc-hostname.mount
[Mount]
Where=/etc/hostname
What=*
Type=noauto
Options=0
