# etc/systemd/system/etc-hosts.mount
[Mount]
Where=/etc/hosts
What=*
Type=noauto
Options=0
