# etc/systemd/system/etc-resolv.conf.mount
[Mount]
Where=/etc/resolv.conf
What=*
Type=noauto
Options=0
